# src/etc/service/freeradius/run
#!/bin/sh
exec /usr/sbin/freeradius -f $FREERADIUS_OPTS
