# src/etc/service/freeradius/run
#!/bin/sh
export KRB5_CLIENT_KTNAME=/etc/freeradius/radius.keytab 
/usr/bin/kdestroy -A
exec /usr/sbin/freeradius -f $FREERADIUS_OPTS
